# cleanMedia/com.nulleffect.clean-movie-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nulleffect.clean-movie-watcher</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/stephenleonard/git/nulleffect/CleanMedia/cleanMedia/watch-movie.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/Volumes/Seagate</key>
            <true/>
        </dict>
    </dict>

    <!-- Backstop against a fast respawn loop (launchd's default floor is 10s). -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- stdout discarded: watch-movie.sh tees everything to logs/watch-movie.log
         (rotated daily). stderr is kept (small). -->
    <key>StandardOutPath</key>
    <string>/dev/null</string>

    <key>StandardErrorPath</key>
    <string>/Users/stephenleonard/git/nulleffect/CleanMedia/logs/launchd-movie-stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/stephenleonard/git/nulleffect/CleanMedia</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
