# cleanMedia/com.nulleffect.clean-tv-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nulleffect.clean-tv-watcher</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/stephenleonard/git/nulleffect/CleanMedia/cleanMedia/watch-tv.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/Volumes/Seagate</key>
            <true/>
        </dict>
    </dict>

    <!-- Backstop: never respawn the script faster than this many seconds, even
         if it exits immediately. launchd's default floor is 10s, which is what
         drove the 2026-06-28 ~11s respawn loop. -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- stdout is discarded: watch-tv.sh already tees everything to
         logs/watch-tv.log (rotated daily). This stops a redundant multi-GB
         launchd-stdout.log from being written. stderr is kept (small). -->
    <key>StandardOutPath</key>
    <string>/dev/null</string>

    <key>StandardErrorPath</key>
    <string>/Users/stephenleonard/git/nulleffect/CleanMedia/logs/launchd-stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/stephenleonard/git/nulleffect/CleanMedia</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
