# cleanMedia/com.nulleffect.clean-log-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nulleffect.clean-log-rotate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/stephenleonard/git/nulleffect/CleanMedia/cleanMedia/rotate-logs.sh</string>
    </array>

    <!-- Run daily at 00:05 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/stephenleonard/git/nulleffect/CleanMedia/logs/rotate.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/stephenleonard/git/nulleffect/CleanMedia/logs/rotate.log</string>
</dict>
</plist>
