# scripts/com.nulleffect.clean-tv-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nulleffect.clean-tv-watcher</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/stephenleonard/git/nulleffect/CleanMedia/scripts/watch-tv.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/Volumes/Seagate</key>
            <true/>
        </dict>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/stephenleonard/git/nulleffect/CleanMedia/logs/launchd-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/stephenleonard/git/nulleffect/CleanMedia/logs/launchd-stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/stephenleonard/git/nulleffect/CleanMedia</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
